# udev/configure-printer@.service.in
[Unit]
Description=Configure Plugged-In Printer
Requires=cups.socket
After=cups.socket

[Service]
ExecStart=@udevhelperdir@/configure-printer add "%i"
ExecStop=@udevhelperdir@/configure-printer remove "%i"
RemainAfterExit=Yes

[Install]
#WantedBy=printer.target
